# service/log/run
#!/bin/execlineb -P

/bin/s6-log T /var/log/vind
