# service/log/run
#!/bin/execlineb -P

/usr/local/bin/s6-log T /var/log/vind
